# roles/ente/templates/podman-ente-postgres.container.j2
[Unit]
Description=Postgres for ente

[Container]
Image=docker.io/library/postgres:{{ ente_postgres_version }}
AutoUpdate=registry

Environment=POSTGRES_USER={{ vault_ente_db_user }}
Environment=POSTGRES_PASSWORD={{ vault_ente_db_pass }}
Environment=POSTGRES_DB={{ ente_db_name }}

Volume=ente.postgres-data.volume:/var/lib/postgresql/data

HealthCmd=pg_isready -q -d {{ ente_db_name }} -U {{ vault_ente_db_user }}
HealthInterval=1s
HealthTimeout=5s
HealthRetries=20

Network=ente.network

[Service]
Restart=always
# Extend Timeout to allow time to pull the image
TimeoutStartSec=600

[Install]
# Start by default on boot
WantedBy=multi-user.target default.target podman-ente-museum.service
